# repotools/templates/promote-nightly.cron.in
# Promote latest builds from testing branch
# Occurs at 10pm Pacific every day
0 6 * * * root /bin/bash @NIGHTLY_PATH@/promote-nightly.sh

